# Kodi.mw
== XBMC Autostart ==

=== xinitrc ===
Both methods require you to create an xbmc user and place an .xinitrc:

First create the xbmc user
{{RootCmd|useradd -m -G audio,cdrom,video,cdrw,usb,users xbmc}}


{{File|/home/xbmc/.xinitrc||<pre>
# .xinitrc
[ -x /usr/bin/nvidia-settings ] && /usr/bin/nvidia-settings -l
/usr/bin/xset s noblank
/usr/bin/xset s off
/usr/bin/xset -dpms
/usr/bin/evilwm &
exec /usr/bin/xbmc
</pre>}}

=== Via inittab ===
This is based on the way the mythtv ebuilds enable the autostart of the mythtv frontend. But since XBMC is the better fronted, I want XBMC to autostart on boot. :)

Then create the login scripts
{{File|/home/xbmc/.bash_profile||<pre>
/home/xbmc/.bash_profile
# .bash_profile
case "`tty`" in
        *tty9) startx; logout ;;
esac
</pre>}}

Finally add an inittab entry which does the autologin
{{File|/etc/inittab||<pre>
...
c9:2345:respawn:/sbin/mingetty --autologin xbmc tty9
...
</pre>}}

=== With init script ===
{{File|/etc/init.d/xbmc||<pre>
#!/sbin/runscript
# Released under the 2-clause BSD license.

description="Start and stop XBMC"

depend()
{
	need net
}



start()
{
	ebegin "Starting XBMC"
	su - xbmc -c "startx" & > /dev/null
	eend 0
}

stop()
{
	killall xbmc.bin
	eend $?
}
</pre>}}

== Shutdown ==
Allow the xbmc user to issue an shutdown via [[polkit]], resulting in a shutdown option in the XBMC menu:
{{File|/etc/polkit-1/rules.d/60-xbmc.rules||<pre>
polkit.addRule(function(action, subject) {
        if (( (action.id.indexOf("org.freedesktop.udisks.") == 0) || (action.id.indexOf("org.freedesktop.upower.")== 0) || (action.id.indexOf("org.freedesktop.consolekit.")== 0) ) &&
        subject.user=="xbmc") {
        return polkit.Result.YES;
        }
});
</pre>}}

Note: ensure you have the '''udev''' USE flag set for shutdown to work

== Hardened profile notes ==

XBMC runs in hardened profiles. One caveat is that to get networking fully working, '''CONFIG_GRKERNSEC_PROC_USER''' kernel option has to be unset. If you want to restrict '''/proc''' access, set '''CONFIG_GRKERNSEC_PROC_USERGROUP''' and '''GRKERNSEC_PROC_GID''' to a special group that will have access to '''/proc''' filesystem. You can then add add xbmc user to this group.

==External resources==
* http://xbmc.org/

[[Category:Multimedia]]
