# Rygel.mw
'''THIS SETUP IS UNSUPPORTED BY GENTOO, DO NOT REPORT BUGS OR COMPLAIN TO THE MAINTAINER IF THINGS DO NOT WORK'''
__TOC__
== Setup Rygel upnp/dlna media server ==
For this howto we can not (yet) use the tracker and dbus interface. The reason is that we start rygel as a daemon with an init script. This is means we do not have a dbus user session and thus can not use tracker.

This howto is aplit over the following sections:
* emerge rygel
* create directories
* create user and group
* copy the files

=== init script ===
{{File|/et/init.d/rygel||<pre>
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
		need net
}

start() {
		ebegin "Starting Rygel"
		start-stop-daemon --start  --background --quiet \
				--make-pidfile \
				--pidfile ${RYGEL_PIDFILE} \
				--user ${RYGEL_USER} \
				--group ${RYGEL_GROUP} \
				--stdout ${RYGEL_LOGFILE} \
				--stderr ${RYGEL_LOGFILE} \
				--exec /usr/bin/rygel -- --config ${RYGEL_CONFIG} \
				${RYGEL_OPTIONS}
		eend $?

}

stop() {
		ebegin "Stopping Rygel"
		start-stop-daemon --stop --retry 10 --quiet --pidfile ${RYGEL_PIDFILE}
		eend $?

}</pre>}}

{{File|/etc/conf.d/rygel||<pre>
# /etc/conf.d/rygel: config file for /etc/init.d/rygel

#Network interface
#RYGEL_IFACE="lo"
#Port
#RYGEL_PORT=12345
#Rygel pid file
RYGEL_PIDFILE="/var/run/rygel.pid"
#Log file
RYGEL_LOGFILE="/var/log/rygel.log"
#Run as user
RYGEL_USER="rygel"
#Rygel group
RYGEL_GROUP="storage"
#Path to config file
RYGEL_CONFIG="/etc/rygel.conf"
#Rygel options
RYGEL_OPTIONS="--disable-transcoding"</pre>}}
